# proiect dizertatie/rc.local
#!/bin/sh -e
#
# rc.local
 
mount /dev/mmcblk0p3 /var/www/capture_files
/root/autoscan > /dev/null 2>&1 &
exit 0